# srt-server.conf
# ============================================================================
# Supervisor Configuration for SRT Server
# 
# This configuration manages the SRT server process that listens on port 9000
# for incoming SRT streams from vMix/OBS encoders and relays them to 
# Flussonic RTMP server.
#
# Install:
#   sudo cp srt-server.conf /etc/supervisor/conf.d/
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl start srt-server
#
# Monitor:
#   sudo supervisorctl status srt-server
#   sudo tail -f /var/www/mediaserver/storage/logs/srt-server.log
# ============================================================================

[program:srt-server]
process_name=%(program_name)s_%(process_num)02d
command=/bin/bash /var/www/mediaserver/srt-server.sh compassiontv
directory=/var/www/mediaserver
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/mediaserver/storage/logs/srt-server.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
user=root
priority=999

# Restart policy
; If process takes more than startsecs to start then we consider it a success
autorestart=true
; Number of seconds the spawned process has to stay running to consider the
; start successful (moving it from STARTING to RUNNING state)
startsecs=10
; Whether the process group will be killed when the process receives a fatal
; signal (SIGSEGV), SIGILL, or SIGABRT)
stopasgroup=true
; Supervisor will kill the entire process group
killasgroup=true
; Priority of the program in the start/shutdown queue
priority=999
