# srt-server.conf
# ============================================================================
# Supervisor Configuration for SRT Server
# 
# This configuration manages the SRT server process that listens on port 9000
# for incoming SRT streams from vMix/OBS encoders and relays them to 
# Flussonic RTMP server.
#
# Install:
#   sudo cp srt-server.conf /etc/supervisor/conf.d/
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl start srt-server
#
# Monitor:
#   sudo supervisorctl status srt-server
#   sudo tail -f /var/www/mediaserver/storage/logs/srt-server.log
# ============================================================================

[program:srt-server]
process_name=%(program_name)s
command=/usr/bin/python3 /var/www/mediaserver/srt-server.py
directory=/var/www/mediaserver
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/www/mediaserver/storage/logs/srt-server.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
user=root
priority=999

